# systemd/piglot.service
[Unit]
Description=PiGlot - AI Language Learning Assistant
After=network-online.target sound.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/piglot
ExecStart=/home/pi/piglot/.venv/bin/piglot
Restart=on-failure
RestartSec=5
Environment=PYTHONUNBUFFERED=1

# Audio permissions
SupplementaryGroups=audio

[Install]
WantedBy=multi-user.target
